# dots/crontab/crontab
MAILTO=kueppo

# change wallpaper every 10
* * * * *	chwal /mnt2/media/images
